# Desktop/tgBOT222igorkapa/cryptobot.service
[Unit]
Description=Crypto Signal Bot - Monitor cryptocurrency market and send buy signals
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/tgBOT222igorkapa
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /root/tgBOT222igorkapa/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cryptobot

# Ограничения ресурсов
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

